# buildroot/board/rockchip/px30/fs-overlay/etc/udev/rules.d/11usb.rules
SUBSYSTEM=="udc",ACTION=="change",DRIVER=="configfs-gadget",RUN+="/usr/bin/usb_config %E{DEVPATH}"
